# src/paperless-office-backend/startupjob.conf
description "startup job for starting services"
author "Joey Driessen"
start on runlevel [2345]
task
exec /home/PaperlessOffice/paperless-office/src/paperless-office-backend/startupScript.sh
